# contrib/screenpipe-autostart/com.screenpipe.recorder.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Screenpipe 开机自启 LaunchAgent 模板
  使用方式见同目录 README.md
  占位符：
    __USER_HOME__ — 替换为 /Users/<你的用户名>
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.screenpipe.recorder</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/screenpipe</string>
        <string>record</string>
        <string>--disable-audio</string>
        <string>-l</string>
        <string>chinese</string>
        <string>--use-pii-removal</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>__USER_HOME__</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>__USER_HOME__</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>Crashed</key>
        <true/>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>60</integer>

    <key>StandardOutPath</key>
    <string>__USER_HOME__/Library/Logs/screenpipe.out.log</string>

    <key>StandardErrorPath</key>
    <string>__USER_HOME__/Library/Logs/screenpipe.err.log</string>

    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>5</integer>
</dict>
</plist>
